# docker_config/dockerfiles/crontab
* * * * * cd /var/www/html/imis && php artisan schedule:run >> /dev/null 2>&1
